# github-inbox-bot.service
[Unit]
Description=GitHub Inbox Discord Bot
After=network.target

[Service]
Type=simple
User=cappo
WorkingDirectory=/home/cappo/github-inbox-bot
Environment=PATH=/home/cappo/github-inbox-bot/venv/bin
ExecStart=/home/cappo/github-inbox-bot/venv/bin/python bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
